# supervisor_collector.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
childlogdir=/var/log/supervisor

# Metasploit JSON RPC server via puma.
# Runs as the dedicated 'msf' user so that msfdb config (~/.msf4/database.yml)
# is readable.  MSF_WS_JSON_RPC_API_TOKEN is read from the token file written
# during installation; puma must start before the collector.
[program:msf_json_rpc]
command=/bin/bash -c "MSF_WS_JSON_RPC_API_TOKEN=$(cat /opt/collector/msf_rpc_token) RACK_ENV=production /opt/metasploit-framework/embedded/bin/ruby /opt/metasploit-framework/embedded/lib/ruby/gems/3.4.0/gems/puma-6.6.0/bin/puma /opt/metasploit-framework/embedded/framework/msf-json-rpc.ru --bind tcp://127.0.0.1:8081"
directory=/opt/metasploit-framework/embedded/framework
autostart=true
autorestart=true
startsecs=60
startretries=3
stderr_logfile=/var/log/supervisor/msf_json_rpc_stderr.log
stdout_logfile=/var/log/supervisor/msf_json_rpc_stdout.log
user=msf
environment=HOME="/home/msf"

[program:collector] 
command=/bin/bash -c "source /root/venv/bin/activate; MSF_JSON_RPC_TOKEN=$(cat /opt/collector/msf_rpc_token) python3 /opt/collector/scan_poller.py -x $(cat /root/.collector_api_key)"
directory=/opt/collector
autostart=true 
autorestart=true 
startsecs=30  
startretries=3  
stderr_logfile=/var/log/supervisor/collector_supervisor_stderr.log
stdout_logfile=/var/log/supervisor/collector_supervisor_stdout.log
user=root

